# scripts/install/betrobot.cron
#!/bin/bash -x


cd /root/betrobot
echo "==================================================" >> update.log
echo $(date -R) >> update.log

echo "--------------------------------------------------" >> update.log
./scripts/update/whoscored/update.sh >> update.log 2>&1

echo "--------------------------------------------------" >> update.log
./scripts/update/betcity/update.sh >> update.log 2>&1

echo "--------------------------------------------------" >> update.log
./scripts/update/intelbet/update.sh >> update.log 2>&1

echo "--------------------------------------------------" >> update.log
./scripts/update/propose/propose.sh >> update.log 2>&1

echo "--------------------------------------------------" >> update.log
# FIXME: Перезагрузить сервер

echo $(date -R) >> update.log
echo "==================================================" >> update.log
echo >> update.log
