# install/betrobot.cron
#!/bin/bash -xe

cd /root/betrobot
echo $(date -R) >> update.log

./betrobot/betcity/update_current.sh >> update.log 2>&1

python3 betrobot/propose_default.sh >> update.log 2>&1

echo $(date -R) >> update.log
echo >> update.log
echo >> update.log
